# backend/supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/tmp/supervisord.pid

[program:redis]
command=redis-server --protected-mode no --port 6379 --bind 127.0.0.1 --dir /tmp
autostart=true
autorestart=true
stdout_logfile=/app/logs/redis.log
stderr_logfile=/app/logs/redis.err.log

[program:ollama]
command=ollama serve
autostart=true
autorestart=true
environment=OLLAMA_HOST="127.0.0.1:11434"
stdout_logfile=/app/logs/ollama.log
stderr_logfile=/app/logs/ollama.err.log

[program:chromadb]
command=chroma run --host 127.0.0.1 --port 8001 --path /app/chroma_data
autostart=true
autorestart=true
stdout_logfile=/app/logs/chromadb.log
stderr_logfile=/app/logs/chromadb.err.log

[program:django]
command=daphne -b 0.0.0.0 -p 7860 config.asgi:application
autostart=true
autorestart=true
environment=REDIS_URL="redis://127.0.0.1:6379/0",CHROMA_HOST="127.0.0.1",CHROMA_PORT="8001",OLLAMA_BASE_URL="http://127.0.0.1:11434"
stdout_logfile=/app/logs/django.log
stderr_logfile=/app/logs/django.err.log

[program:celery-worker]
command=celery -A config worker -l info -c 1 --max-tasks-per-child=2 -Q default,analysis,ai
autostart=true
autorestart=true
environment=REDIS_URL="redis://127.0.0.1:6379/0",CHROMA_HOST="127.0.0.1",CHROMA_PORT="8001",OLLAMA_BASE_URL="http://127.0.0.1:11434"
stdout_logfile=/app/logs/celery-worker.log
stderr_logfile=/app/logs/celery-worker.err.log

[program:celery-beat]
command=celery -A config beat -l info --scheduler django_celery_beat.schedulers:DatabaseScheduler
autostart=true
autorestart=true
environment=REDIS_URL="redis://127.0.0.1:6379/0"
stdout_logfile=/app/logs/celery-beat.log
stderr_logfile=/app/logs/celery-beat.err.log

[program:onchain-monitor]
command=python manage.py run_onchain_monitor
autostart=true
autorestart=true
environment=REDIS_URL="redis://127.0.0.1:6379/0"
stdout_logfile=/app/logs/onchain-monitor.log
stderr_logfile=/app/logs/onchain-monitor.err.log

[program:ingest-knowledge]
command=bash -c "sleep 10 && python manage.py ingest_knowledge --source all"
autostart=true
autorestart=false
startsecs=0
environment=REDIS_URL="redis://127.0.0.1:6379/0",CHROMA_HOST="127.0.0.1",CHROMA_PORT="8001",OLLAMA_BASE_URL="http://127.0.0.1:11434"
stdout_logfile=/app/logs/ingest-knowledge.log
stderr_logfile=/app/logs/ingest-knowledge.err.log
